# setup/nepse-trading.timer
[Unit]
Description=NEPSE AI — Trading Loop Timer (every 6 min, 11:00 AM–3:00 PM NST)

[Timer]
# Market hours: 11:00 AM–3:00 PM NST = 05:15–09:15 UTC
# Run every 15 minutes within that window.
# calendar_guard inside main.py will exit cleanly outside market hours —
# so it's safe to run this timer slightly wider (05:15–09:18 UTC).
OnCalendar=Mon,Tue,Wed,Thu,Fri 05:15:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 05:30:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 05:45:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 06:00:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 06:15:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 06:30:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 06:45:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 07:00:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 07:15:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 07:30:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 07:45:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 08:00:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 08:15:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 08:30:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 08:45:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 09:00:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 09:15:00 UTC

Persistent=false
AccuracySec=10s

[Install]
WantedBy=timers.target
