# setup/nepse-sleep.timer
[Unit]
Description=NEPSE AI — Sleep Scheduler Timer (April 2026)

[Timer]
OnCalendar=Mon,Tue,Wed,Thu,Fri 10:00:00 UTC
OnCalendar=Sun 12:30:00 UTC
Persistent=false
AccuracySec=60s

[Install]
WantedBy=timers.target