# setup/nepse-summary.timer
[Unit]
Description=NEPSE AI — Nightly Summary Timer (3:30 pm PM NST = 09:30 UTC, Mon–Fri)

[Timer]
# 3:30 PM NST = 9:30 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 9:45:00 UTC
Persistent=true
AccuracySec=60s

[Install]
WantedBy=timers.target
