# setup/nepse-morning.timer
[Unit]
Description=NEPSE AI — Morning Workflow Timer (10:45 AM NST = 05:00 UTC, Mon–Fri)

[Timer]
# Updated April 2026: Monday-Friday
OnCalendar=Mon,Tue,Wed,Thu,Fri 05:00:00 UTC
Persistent=true
AccuracySec=30s

[Install]
WantedBy=timers.target
