# setup/nepse-trading.timer
[Unit]
Description=NEPSE AI — Trading Loop Timer (every 30 min, 11:15 AM–2:45 PM NST)

[Timer]
# Market hours window: 11:15 AM–2:45 PM NST = 05:30–09:00 UTC
# Run every 30 minutes within that window.
# calendar_guard inside main.py will exit cleanly outside market hours.
OnCalendar=Mon,Tue,Wed,Thu,Fri 05:30:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 06:00:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 06:30:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 07:00:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 07:30:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 08:00:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 08:30:00 UTC
OnCalendar=Mon,Tue,Wed,Thu,Fri 09:00:00 UTC

Persistent=false
AccuracySec=10s

[Install]
WantedBy=timers.target