# setup/nepse-gitpull.timer
[Unit]
Description=NEPSE AI — Git Pull Timer

[Timer]
OnCalendar=Mon,Tue,Wed,Thu,Fri 04:25:00 UTC
Persistent=false
AccuracySec=60s

[Install]
WantedBy=timers.target